# challenges/in_the_zone/deployment/supervisord.conf
[supervisord]
nodaemon=true

[program:nsd]
command=nsd -d

[program:nginx]
command=nginx -g "daemon off;"
